# mirai-knowledge-app.service
[Unit]
Description=Mirai Knowledge Systems Flask Application (Production)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=kensan
WorkingDirectory=/mnt/LinuxHDD/Mirai-Knowledge-Systems/backend

# 本番環境ポート: HTTP 9100, HTTPS 9443
ExecStart=/mnt/LinuxHDD/Mirai-Knowledge-Systems/venv_linux/bin/gunicorn \
    --config /mnt/LinuxHDD/Mirai-Knowledge-Systems/backend/gunicorn.conf.py \
    app_v2:app

Restart=always
RestartSec=5

# 本番環境設定
Environment=MKS_ENV=production
Environment=MKS_HTTP_PORT=9100
Environment=MKS_HTTPS_PORT=9443
Environment=MKS_DEBUG=false
Environment=MKS_USE_POSTGRESQL=true
Environment=MKS_LOAD_SAMPLE_DATA=false
Environment=MKS_CREATE_DEMO_USERS=false
Environment=MKS_FORCE_HTTPS=true
Environment=MKS_SECRET_KEY=prod-secret-key-change-me-2026
Environment=MKS_JWT_SECRET_KEY=prod-jwt-secret-key-change-me-2026
Environment=MKS_CORS_ORIGINS=https://192.168.0.187:9443,https://192.168.0.145:9443

[Install]
WantedBy=multi-user.target