# config/mirai-knowledge-system-dev.service
[Unit]
Description=Mirai Knowledge System - 建設土木ナレッジシステム (Development)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=kensan
Group=kensan
WorkingDirectory=/mnt/LinuxHDD/Mirai-Knowledge-Systems/backend
EnvironmentFile=/mnt/LinuxHDD/Mirai-Knowledge-Systems/backend/.env
Environment="PATH=/mnt/LinuxHDD/Mirai-Knowledge-Systems/backend/venv_linux/bin:/usr/local/bin:/usr/bin:/bin"
Environment="MKS_ENV=development"
Environment="MKS_DEBUG=false"
# Python実行ファイルのフルパス指定
ExecStart=/usr/bin/python3 /mnt/LinuxHDD/Mirai-Knowledge-Systems/backend/app_v2.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mirai-knowledge-dev

# セキュリティ設定
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
